# dotfiles/profile.d/dircolors.sh
# ~/.profile/dircolors.sh

if ihave dircolors; then
    if [[ -r ~/.dircolors ]]; then
        eval "$(dircolors -b ~/.dircolors)"
    else
        eval "$(dircolors -b)"
    fi
fi
